# package/extra/luci-app-udp2raw/root/etc/hotplug.d/iface/99-udp2raw
#!/bin/sh

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] &&
	! /etc/init.d/udp2raw running instance1 &&
	/etc/init.d/udp2raw reload
