# package/extra/luci-app-udp2raw/root/etc/hotplug.d/iface/99-udp2raw
#!/bin/sh

[ "$ACTION" = "ifup" ] &&
	[ "$INTERFACE" = "wan" ] &&
	[ -e "/var/run/udp2raw.hostname_not_resolved.err" ] && {
	/etc/init.d/udp2raw reload
}
